# storm.upstart.conf
# Ubuntu upstart file at /etc/init/storm.conf

description "Storm"

limit nofile 32768 32768

start on runlevel [2345]
stop on [!12345]

respawn
respawn limit 2 5

env CONFIG_HOME=/opt/storm/config
env STORM_HOME=/opt/storm
umask 007

kill timeout 300

setuid storm
setgid storm

script
  . /etc/default/storm

  if [ "x$ENABLE" = "xyes" ]; then
    exec $STORM_HOME/bin/kafka-server-start.sh $CONFIG_HOME/server.properties
  fi
end script
